# ops/installfiles/etc/systemd/system/mongod.service.d/override.conf
[Service]
ExecStartPost=/bin/sh -c 'while ! /usr/bin/mongosh --eval "db.version()" > /dev/null 2>&1; do sleep 0.1; done'
